# net-p2p/eclair-bin/files/eclair_at.service
[Unit]
Description=Eclair lightning network node
After=network.target

[Service]
ExecStart=/usr/bin/eclair-node
User=eclair-%i
