# net-p2p/eclair-bin/files/eclair_at.service
[Unit]
Description=Eclair lightning network node
After=network.target

[Service]
ExecStart=/usr/bin/eclair-node \
			-Declair.datadir=/var/lib/eclair-%i \
			-Declair.printToConsole

User=eclair-%i

[Install]
WantedBy=multi-user.target
